# templates/default/datadisk.conf.erb
start on local-filesystems

task

script
  rsync -aH <%= @data_part %>/ <%= @data_part %>_temp_mountpoint/
  mv <%= @data_part %> <%= @data_part %>_old
  mkdir <%= @data_part %>
  umount <%= @data_part %>_temp_mountpoint
  sed -i -e 's:<%= @data_part %>_temp_mountpoint:<%= @data_part %>:' /etc/fstab
  mount <%= @data_part %>
  sed -i -e 's/DEFAULT_RUNLEVEL=1/DEFAULT_RUNLEVEL=2/' /etc/init/rc-sysinit.conf
  # harakiri
  rm -f /etc/init/datadisk.conf
end script

post-stop script
  init 2
end script
